# systemd/vpn-api.service
[Unit]
Description=VPN Key Management API
After=network.target
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=vpnapi
Group=vpnapi
WorkingDirectory=/root/vpn-server
Environment=PATH=/root/vpn-server/venv/bin
EnvironmentFile=/root/vpn-server/.env
ExecStart=/bin/bash -c 'if [ "${VPN_ENABLE_HTTPS:-false}" = "true" ] && [ -f "${VPN_SSL_CERT_PATH:-/etc/ssl/certs/vpn-api.crt}" ] && [ -f "${VPN_SSL_KEY_PATH:-/etc/ssl/private/vpn-api.key}" ]; then /root/vpn-server/venv/bin/uvicorn api:app --host ${VPN_HOST:-0.0.0.0} --port ${VPN_PORT:-8000} --workers ${VPN_WORKERS:-2} --limit-max-requests ${VPN_WORKER_MAX_REQUESTS:-500} --ssl-keyfile ${VPN_SSL_KEY_PATH} --ssl-certfile ${VPN_SSL_CERT_PATH}; else /root/vpn-server/venv/bin/uvicorn api:app --host ${VPN_HOST:-0.0.0.0} --port ${VPN_PORT:-8000} --workers ${VPN_WORKERS:-2} --limit-max-requests ${VPN_WORKER_MAX_REQUESTS:-500}; fi'
Restart=always
RestartSec=5
LimitNOFILE=65535
MemoryMax=512M
UMask=0077
NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
RestrictNamespaces=yes
RestrictRealtime=yes
LockPersonality=yes
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
ReadWritePaths=/root/vpn-server /var/log /run

[Install]
WantedBy=multi-user.target